# scripts/systemd/atlas-autopush.service
[Unit]
Description=ATLAS-PTQ auto-snapshot + git push
After=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /capstor/scratch/cscs/fnoorala/ATLAS-PTQ/scripts/auto_push.sh
